# hybrid_build/workspace/rootfs/lib/apk/packages/base-files.list
/bin/board_detect
/bin/config_generate
/bin/ipcalc.sh
/etc/TZ
/etc/apk/repositories.d/distfeeds.list
/etc/banner
/etc/banner.failsafe
/etc/board.d/01_led
/etc/board.d/02_network
/etc/board.d/03_gpio_switches
/etc/board.d/99-default_network
/etc/device_info
/etc/diag.sh
/etc/ethers
/etc/fstab
/etc/group
/etc/hosts
/etc/hotplug.d/leds/00-init
/etc/hotplug.d/net/00-sysctl
/etc/init.d/boot
/etc/init.d/done
/etc/init.d/gpio_switch
/etc/init.d/led
/etc/init.d/sysctl
/etc/init.d/sysfixtime
/etc/init.d/system
/etc/init.d/umount
/etc/inittab
/etc/iproute2/ematch_map
/etc/iproute2/rt_protos
/etc/iproute2/rt_tables
/etc/localtime
/etc/mtab
/etc/openwrt_release
/etc/openwrt_version
/etc/os-release
/etc/passwd
/etc/preinit
/etc/profile
/etc/profile.d/00-passwordless-root.sh
/etc/profile.d/apk-cheatsheet.sh
/etc/protocols
/etc/rc.button/failsafe
/etc/rc.button/power
/etc/rc.button/reboot
/etc/rc.button/reset
/etc/rc.button/rfkill
/etc/rc.common
/etc/rc.local
/etc/resolv.conf
/etc/services
/etc/shadow
/etc/shells
/etc/shinit
/etc/sysctl.conf
/etc/sysctl.d/10-default.conf
/etc/sysupgrade.conf
/etc/uci-defaults/05-migrate-ten64-gpio
/etc/uci-defaults/10_migrate-shadow
/etc/uci-defaults/11_network-migrate-bridges
/etc/uci-defaults/12_network-generate-ula
/etc/uci-defaults/14_network-generate-clientid
/etc/uci-defaults/50-root-passwd
/lib/functions.sh
/lib/functions/caldata.sh
/lib/functions/ipv4.sh
/lib/functions/leds.sh
/lib/functions/migrations.sh
/lib/functions/network.sh
/lib/functions/preinit.sh
/lib/functions/service.sh
/lib/functions/system.sh
/lib/functions/uci-defaults.sh
/lib/preinit/00_preinit.conf
/lib/preinit/01_sysinfo_acpi
/lib/preinit/02_default_set_state
/lib/preinit/02_sysinfo
/lib/preinit/10_indicate_failsafe
/lib/preinit/10_indicate_preinit
/lib/preinit/30_failsafe_wait
/lib/preinit/40_run_failsafe_hook
/lib/preinit/50_indicate_regular_preinit
/lib/preinit/70_initramfs_test
/lib/preinit/79_move_config
/lib/preinit/80_mount_root
/lib/preinit/99_10_failsafe_login
/lib/preinit/99_10_run_init
/lib/upgrade/common.sh
/lib/upgrade/do_stage2
/lib/upgrade/fwtool.sh
/lib/upgrade/keep.d/base-files
/lib/upgrade/keep.d/base-files-essential
/lib/upgrade/platform.sh
/lib/upgrade/stage2
/lib/upgrade/tar.sh
/rom/note
/sbin/firstboot
/sbin/hotplug-call
/sbin/led.sh
/sbin/sysupgrade
/usr/lib/os-release
/usr/libexec/login.sh
/usr/libexec/validate_firmware_image
/var
