# hybrid_build/workspace/rootfs/lib/apk/packages/logd.list
/etc/init.d/log
/lib/apk/packages/logd.alternatives
/lib/apk/packages/logd.rusers
/sbin/logd
/usr/libexec/logread-ubox
/usr/share/acl.d/logd.json
