# hybrid_build/workspace/rootfs/lib/apk/packages/netifd.list
/etc/hotplug.d/iface/00-netstate
/etc/init.d/network
/etc/init.d/packet_steering
/etc/uci-defaults/14_migrate-dhcp-release
/etc/udhcpc.user
/lib/netifd/dhcp.script
/lib/netifd/main.uc
/lib/netifd/netifd-proto.sh
/lib/netifd/proto/dhcp.sh
/lib/netifd/utils.sh
/lib/netifd/utils.uc
/lib/network/config.sh
/lib/upgrade/keep.d/netifd
/sbin/devstatus
/sbin/ifdown
/sbin/ifstatus
/sbin/ifup
/sbin/netifd
/usr/libexec/network/packet-steering.uc
/usr/share/udhcpc/default.script
